# render-supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log

[program:gateway]
command=openclaw gateway run
directory=/root
autostart=true
autorestart=true
stdout_logfile=/tmp/gateway.log
stderr_logfile=/tmp/gateway-error.log

[program:backend]
command=uvicorn backend.server:app --host 0.0.0.0 --port %(ENV_PORT)s --workers 1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/tmp/backend.log
stderr_logfile=/tmp/backend-error.log
